# package/p2pmpi/usr/share/p2pmpi/p2pmpi.conf
# p2pmpi - Starts the p2pmpi service. 
#

description	"Start p2pmpi service"

start on net-device-up

script
	/usr/bin/mpiboot
end script
